# docker/etc/s6/polipo/run
#!/bin/bash

# 配置
source ./setup

# 启动代理转换服务
exec su-exec "${USERNAME}" /opt/polipo/polipo -c /etc/polipo/config
