# docker/etc/s6/fastgithub/run
#!/bin/bash

# 启动FastGithub服务
exec su-exec "${USERNAME}" /opt/fastgithub/fastgithub
