# docker/etc/s6/fastgithub/run
#!/bin/bash

# 启动FastGithub服务
/opt/fastgithub/fastgithub
